# docker/supervisord/supervisord.conf
[supervisord]
nodaemon=true
user=root
childlogdir=/var/www/storage/logs/
logfile=/var/www/storage/logs/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB                         ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10                            ; (num of main logfile rotation backups;default 10)
loglevel=debug                                ; (log level;default info; others: debug,warn,trace)
pidfile=/var/www/storage/logs/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false                                ; (start in foreground if true;default false)
minfds=1024                                   ; (min. avail startup file descriptors;default 1024)
minprocs=200                                  ; (min. avail process descriptors;default 200)

[program:swoole]
directory=/var/www
command=sh -c "while true; do php artisan swoole:http start; done"
autostart=true

[program:watcher]
directory=/var/www
command=sh -c "yarn watcher"
autostart=true